# bash/crontab
crontab -e

0 17 * * 4  cd /root/nps_gen; bash main.sh; cd -
